# nginx-phpfpm/etc/supervisor.d/php.conf
[program:php]
command = /usr/sbin/php-fpm7
user = root
autostart = true